# tblog.conf
[program:tblog]
command = gunicorn -w 2 app:app -b localhost:8000
directory = /tblog
