# tblog.conf
[supervisord]
nodaemon=true

[program:tblog]
command = gunicorn -w 2 app:app -b 0.0.0.0:8000
directory = /tblog
